# src/min-cloud-agent.service.in
[Unit]
Description=Process OpenStack/AWS metadata and userdata
ConditionPathExists=!/var/lib/min-cloud-agent/done.stamp

[Service]
Type=oneshot
ExecStart=@libexecdir@/min-cloud-agent

[Install]
WantedBy=multi-user.target
