# src/min-metadata.service.in
[Unit]
Description=Process OpenStack/AWS metadata

[Service]
Type=oneshot
ExecStart=@libexecdir@/min-metadata-service

[Install]
WantedBy=multi-user.target
