# src/min-metadata.service.in
[Unit]
Description=Process OpenStack/AWS metadata

[Service]
Type=oneshot
Before=multi-user.target sshd.service
ExecStart=@libexecdir@/min-metadata-service
