# src/min-metadata.service.in
[Unit]
Description=Process OpenStack/AWS metadata
ConditionPathExists=!/root/.ssh/authorized_keys
ConditionPathExists=!/var/initial-setup-accounts/done

[Service]
Type=oneshot
ExecStart=@libexecdir@/min-metadata-service

[Install]
WantedBy=multi-user.target
